# DEBIAN/postinst
#!/bin/bash
set -e

echo "[INFO] Installing raspi-server..."

# Install Python dependencies
if [ -f /opt/raspi-server/requirements.txt ]; then
    echo "[INFO] Installing Python requirements..."
    pip3 install -r /opt/raspi-server/requirements.txt
fi

# Reload systemd
systemctl daemon-reload

# Enable service
systemctl enable raspi-server.service

# Start service
systemctl start raspi-server.service

echo "[INFO] raspi-server installed and running."
