# DEBIAN/postinst
#!/bin/bash
set -e

echo "[INFO] Installing raspi-server..."

APP_DIR="/opt/raspi-server"
VENV_DIR="$APP_DIR/venv"

# Create virtual environment if not exists
if [ ! -d "$VENV_DIR" ]; then
    echo "[INFO] Creating virtual environment..."
    python3 -m venv "$VENV_DIR"
fi

# Activate venv and install requirements
if [ -f "$APP_DIR/requirements.txt" ]; then
    echo "[INFO] Installing Python requirements inside virtualenv..."
    source "$VENV_DIR/bin/activate"
    pip install --upgrade pip
    pip install -r "$APP_DIR/requirements.txt"
    deactivate
fi

# Reload systemd
systemctl daemon-reload

# Enable service
systemctl enable raspi-server.service

# Start service
systemctl start raspi-server.service

echo "[INFO] raspi-server installed and running inside virtual environment."
