# DEBIAN/prerm
#!/bin/sh -e
# prerm - runs before package removal/upgrade

SERVICE="raspi-server.service"

case "$1" in
  remove|upgrade)
    if systemctl is-active --quiet "$SERVICE"; then
      systemctl stop "$SERVICE" || true
    fi
    systemctl disable "$SERVICE" || true
    ;;
esac

exit 0
