# etc/systemd/system/raspi-server.service
[Unit]
Description=Raspberry Pi Flask Server
After=network.target

[Service]
Type=simple
User=raspiserver
Group=raspiserver
WorkingDirectory=/opt/raspi-server
ExecStart=/usr/bin/python3 /opt/raspi-server/app.py
Restart=on-failure
RestartSec=5
StandardOutput=syslog
StandardError=syslog
SyslogIdentifier=raspi-server

[Install]
WantedBy=multi-user.target
